# services/jukebox-backend.service
[Unit]
Description=Jukebox API Server
After=network.target sound.target

[Service]
Type=simple
User=jukebox
Group=audio
WorkingDirectory=/home/jukebox/jukebox
Environment=PATH=/home/jukebox/jukebox/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/home/jukebox/jukebox
Environment=VIRTUAL_ENV=/home/jukebox/jukebox/venv
Environment=PULSE_SERVER=unix:/run/user/%i/pulse/native
Environment=XDG_RUNTIME_DIR=/run/user/%i
ExecStart=/home/jukebox/jukebox/venv/bin/python -m uvicorn src.api_server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

# Audio-specific settings
SupplementaryGroups=audio

[Install]
WantedBy=multi-user.target